# templates/consul.freebsd.erb
#!/bin/sh

# $FreeBSD: branches/2016Q4/sysutils/consul/files/consul.in 420942 2016-08-26 19:46:24Z allanjude $
#
# PROVIDE: consul
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# consul_enable (bool): Set to NO by default.
#                               Set it to YES to enable consul.
# consul_user (user):        Set user to run consul.
#                               Default is "consul".
# consul_group (group):      Set group to run consul.
#                               Default is "consul".
# consul_dir (dir):          Set dir to run consul in.
#                               Default is "/var/tmp/consul".
# consul_env (dir):          Set environment variables used with consul
#                               Default is "".

. /etc/rc.subr

name=consul
rcvar=consul_enable
extra_commands="reload"

load_rc_config $name

: ${consul_enable:="NO"}
: ${consul_user:="<%= scope.lookupvar('consul::user') %>"}
: ${consul_group:="consul"}
: ${consul_dir:="/var/tmp/consul"}
: ${consul_env:=""}
: ${consul_args:="-rpc-addr=<%= scope.lookupvar('consul::rpc_addr') %>:<%= scope.lookupvar('consul::rpc_port')%> <%= scope.lookupvar('consul::extra_options') %>"}

pidfile=/var/run/consul.pid
procname="<%= scope.lookupvar('consul::bin_dir') %>/consul"
command="/usr/sbin/daemon"
command_args="-f -p ${pidfile} /usr/bin/env ${consul_env} ${procname} agent -ui-dir=/usr/local/share/consul -config-dir=<%= scope.lookupvar('consul::config_dir') %> ${consul_args}"

start_precmd=consul_startprecmd

consul_startprecmd()
{
        if [ ! -e ${pidfile} ]; then
                install -o ${consul_user} -g ${consul_group} /dev/null ${pidfile};
        fi

        if [ ! -d ${consul_dir} ]; then
                install -d -o ${consul_user} -g ${consul_group} ${consul_dir}
        fi

}

run_rc_command "$1"

